# cron.d/fix-aborted
*/5 * * * * root /usr/local/bin/python3.5 -m seisma fix_aborted >> /var/log/cron.log 2>&1
