# cron.d/rotate-data-in-database
00 5 * * * root /usr/local/bin/python3.5 -m seisma rotate >> /var/log/cron.log 2>&1
